# host/systemd/daemon.service
[Unit]
Description=Millennium Tools Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Run the binary directly (no shell needed)
ExecStart=/home/matzen/millennium/host/daemon --config /home/matzen/millennium/host/daemon.conf
# Optional: wait for network before starting (use absolute path + correct bash)
ExecStartPre=/bin/bash -c 'until /bin/ping -c1 -W1 8.8.8.8 >/dev/null; do sleep 1; done'
Restart=on-failure
RestartSec=2
WorkingDirectory=/home/matzen/millennium/host

[Install]
WantedBy=default.target

